# devops/example.ini
[program:server]
environment=CONSUL_ADDR=127.0.0.1:8500
directory = /root/servers/bin
command = /root/servers/bin/server -log_dir=/root/logs/server -remote-config=server-prod.toml
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 20
stdout_logfile = /root/logs/server/stdout.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 20
stderr_logfile = /root/logs/server/stderr.log